# companion/com.deepxanshu.kino-companion.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.deepxanshu.kino-companion</string>
  <key>ProgramArguments</key>
  <array>
    <string>/Users/deepanshu/devs-workspace/fo/magicstick/companion/.venv/bin/python</string>
    <string>/Users/deepanshu/devs-workspace/fo/magicstick/companion/codex_agents.py</string>
    <string>--serial</string>
  </array>
  <key>WorkingDirectory</key>
  <string>/Users/deepanshu/devs-workspace/fo/magicstick</string>
  <key>RunAtLoad</key>
  <true/>
  <key>KeepAlive</key>
  <true/>
  <key>StandardOutPath</key>
  <string>/tmp/kino-companion.log</string>
  <key>StandardErrorPath</key>
  <string>/tmp/kino-companion.log</string>
</dict>
</plist>
